# backup@.timer
[Unit]
Description=%I backup timer

[Timer]
OnCalendar=daily
AccuracySec=1day

[Install]
WantedBy=timers.target
